# package/base-files/files/etc/profile.d/aliases.sh
#!/bin/sh

# $Id$

alias s='cd ..'
alias p='cd -'
alias du='du -h'
alias free='free -m'
alias df='df -h'
alias ce="cd /etc"
alias c~="cd ~"
alias cd~="cd ~"
alias ch="cd ~"
alias cm="cd /mnt"
alias cmn="cd /mnt/nfs"
alias cr="cd /"
alias nst='netstat -aptun'
alias lsu='ls -l /dev/ttyUSB*'
alias lsi='lsof -i'
